# systemd/rbatch@all.timer
[Unit]
Description=Process rbatch-all nightly

[Timer]
OnCalendar=02:00

[Install]
WantedBy=timers.target
