# systemd/rclone-unit@.service
[Unit]
Description=Transfer %I using rclone (%p)
After=network.target

[Service]
ExecStart=/usr/local/bin/rclone-unit.sh %p "%I"
