# systemd/rclone-all.timer
[Unit]
Description=Process rclone nightly

[Timer]
OnCalendar=02:00

[Install]
WantedBy=timers.target
